# examples/systemd/joplin-md-sync-mcp.service
[Unit]
Description=Joplin Markdown Sync MCP server
Documentation=https://github.com/kogeler/joplin-md-sync
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/joplin-md-sync mcp serve --token-file %h/.config/joplin-md-sync/joplin-token --auth-token-file %h/.config/joplin-md-sync/mcp-token
Restart=on-failure
RestartSec=5s
Environment=PYTHONUNBUFFERED=1

# The daemon only needs read access to its executable and token files.
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectControlGroups=true
ProtectHome=read-only
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictSUIDSGID=true

[Install]
WantedBy=default.target
