# scripts/joplin_terminal_service/systemd/joplin-md-sync.service
[Unit]
Description=Joplin Markdown Sync MCP and GPT Actions service
Documentation=https://github.com/kogeler/joplin-md-sync
Requires=joplin-terminal.service
After=joplin-terminal.service

[Service]
Type=simple
ExecStart={exec_start}
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
KillMode=mixed
UMask=0077

NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectControlGroups=true
ProtectHome=read-only
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictSUIDSGID=true

[Install]
WantedBy=default.target
